# etc/profile.d/export_idf_path.sh
export IDF_PATH="C:/101B_FlashStation/home/Administrator/esp-idf-v3.0"