# etc/profile.d/export_idf_path.sh
export IDF_PATH="C:/msys32/home/Administrator/esp-idf-v3.0"