# etc/profile.d/export_idf_path.sh
export IDF_PATH="C:/Users/$USER/Documents/101B_FlashStationXXX/home/Administrator/esp-idf-v3.0"