# etc/profile.d/export_idf_path.sh
export IDF_PATH="C:/Users/$USER/Documents/101B_FlashStation/home/$USER/esp-idf-v3.0"